# container-files/etc/supervisor.d/httpd.conf
# - Apache en segundo plano

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
