# container-files/etc/supervisor.d/sshd.conf
# - SSH en segundo plano
[program:sshd]
command=/usr/sbin/sshd -D -f /etc/ssh/sshd_config
